# rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "Kids ChatBot Server stopped"
