# ha-entity-resolver/rootfs/etc/services.d/ha-entity-resolver/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish script for HA Entity Resolver
# ==============================================================================

bashio::log.info "HA Entity Resolver stopped"
