# ha-entity-endpoint/rootfs/etc/services.d/ha-entity-endpoint/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle service shutdown
# ==============================================================================

bashio::log.info "HA Entity Endpoint service stopped"
