# ha-entity-endpoint/rootfs/etc/services.d/ha-entity-endpoint/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start HA Entity Endpoint service
# ==============================================================================

bashio::log.info "Starting HA Entity Endpoint..."

# Export environment variables
export ACCESS_TOKEN=$(bashio::config 'access_token')
export HA_URL=$(bashio::config 'ha_url')

# Run the Python application
exec python3 /app/run.py
