# chatbot-websocket/rootfs/etc/services.d/chatbot/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# School Chatbot WebSocket Server - S6 Finish Script
# ==============================================================================

# Log the exit code
bashio::log.info "Service exited with code: ${1}"

# If exit code is not 0, log error
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.error "Service crashed with code ${1}"
fi

exit 0
