# chatbot-websocket/rootfs/etc/services.d/chatbot/run
#!/usr/bin/with-contenv bashio

# Get config from options
export AI_PROVIDER=$(bashio::config 'ai_provider')
export AI_MODEL=$(bashio::config 'ai_model')
export TTS_VOICE_VI=$(bashio::config 'tts_voice_vi')
export TTS_VOICE_EN=$(bashio::config 'tts_voice_en')
export SYSTEM_PROMPT=$(bashio::config 'system_prompt')
export MAX_CONTEXT_MESSAGES=$(bashio::config 'max_context_messages')
export TEMPERATURE=$(bashio::config 'temperature')
export MAX_TOKENS=$(bashio::config 'max_tokens')
export LOG_LEVEL=$(bashio::config 'log_level')

# API Keys
export OPENAI_API_KEY=$(bashio::config 'openai_api_key')
export OPENAI_BASE_URL=$(bashio::config 'openai_base_url')
export DEEPSEEK_API_KEY=$(bashio::config 'deepseek_api_key')

# Display config
bashio::log.info "🚀 Starting School Chatbot WebSocket Server..."
bashio::log.info "📋 Configuration:"
bashio::log.info "   AI Provider: ${AI_PROVIDER}"
bashio::log.info "   AI Model: ${AI_MODEL}"
bashio::log.info "   TTS Voice (VI): ${TTS_VOICE_VI}"
bashio::log.info "   TTS Voice (EN): ${TTS_VOICE_EN}"
bashio::log.info "   Max Context: ${MAX_CONTEXT_MESSAGES}"

bashio::log.info "🌐 Starting server on port 5000..."

# Run application
cd /app
exec python3 -m app.main
