# kids-chatbot-server/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Kids ChatBot Server..."

# Get configuration from Home Assistant
export OPENAI_API_KEY=$(bashio::config 'openai_api_key')
export OPENAI_MODEL=$(bashio::config 'openai_model')
export OPENAI_VOICE=$(bashio::config 'openai_voice')
export OPENAI_LANGUAGE=$(bashio::config 'openai_language')
export PORT=$(bashio::config 'port')
export MAX_TOKENS=$(bashio::config 'max_tokens')
export TEMPERATURE=$(bashio::config 'temperature')
export LOG_LEVEL="INFO"

# Log configuration
bashio::log.info "Configuration loaded successfully"
bashio::log.info "Model: ${OPENAI_MODEL} | Language: ${OPENAI_LANGUAGE} | Voice: ${OPENAI_VOICE}"

# Validate API key
if bashio::var.is_empty "${OPENAI_API_KEY}"; then
    bashio::log.error "OpenAI API key is not configured!"
    bashio::log.error "Please set 'openai_api_key' in the add-on configuration"
    exit 1
fi

# Change to app directory
cd /usr/bin || exit 1

# Start Flask app
exec python3 app.py
