# ha-entity-resolver/rootfs/etc/services.d/resolver/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Resolver service stopped"
