# ha-entity-resolver/rootfs/etc/services.d/resolver/finish
#!/bin/sh
echo "[INFO] Resolver stopped"
