# ha-entity-resolver/rootfs/etc/services.d/resolver/finish
#!/bin/bash
echo "[INFO] Resolver service stopped"