# ha-entity-resolver/rootfs/etc/services.d/resolver/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Finish script for resolver service
# ==============================================================================

s6-svc -O .
