# ha-entity-resolver/rootfs/etc/services.d/resolver/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid abc
cd /app
python3 run.py