# ha-entity-resolver/rootfs/etc/services.d/resolver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# HA Entity Resolver - Startup Script
# ==============================================================================

bashio::log.info "Starting HA Entity Resolver..."

export LOG_LEVEL=$(bashio::config 'log_level')
export HA_URL=$(bashio::config 'ha_url')
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"

# Write access tokens
bashio::config 'access_tokens' | jq -c . > /tmp/access_tokens.json

# Write MQTT prefix
bashio::config 'mqtt_device_prefix' > /tmp/mqtt_prefix.txt

bashio::log.info "Configuration loaded"
bashio::log.info "  HA URL: ${HA_URL}"
bashio::log.info "  MQTT Prefix: $(cat /tmp/mqtt_prefix.txt)"

cd /app || bashio::exit.nok "Cannot change to /app directory"

bashio::log.info "Starting Flask server on port 5003..."

exec python3 run.py
