# ha-entity-resolver/rootfs/etc/services.d/resolver/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting HA Entity Resolver..."

export LOG_LEVEL=$(bashio::config 'log_level')
export HA_URL=$(bashio::config 'ha_url')
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"

bashio::config 'access_tokens' | jq -c . > /tmp/access_tokens.json

if bashio::config.has_value 'mqtt_device_prefix'; then
    bashio::config 'mqtt_device_prefix' > /tmp/mqtt_prefix.txt
else
    echo "192168100131" > /tmp/mqtt_prefix.txt
fi

cd /app
exec python3 run.py
